# systemd/clutch-server.service.template
[Unit]
Description=OpenClutch Next.js Server
Documentation=https://github.com/dbachelder/clutch
After=network.target

[Service]
Type=simple
WorkingDirectory={{WORKING_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart={{NODE_BIN}} ./node_modules/next/dist/bin/next start -p 3002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clutch-server

[Install]
WantedBy=default.target
