# systemd/clutch-session-watcher.service.template
[Unit]
Description=OpenClutch Session Watcher
Documentation=https://github.com/dbachelder/clutch
After=clutch-server.service
Wants=clutch-server.service

[Service]
Type=simple
WorkingDirectory={{WORKING_DIR}}
EnvironmentFile={{ENV_FILE}}
ExecStart={{NODE_BIN}} ./node_modules/tsx/dist/cli.cjs worker/session-watcher.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clutch-session-watcher

[Install]
WantedBy=default.target
