# systemd/trap-server.service
[Unit]
Description=Trap Next.js Server
Documentation=https://github.com/dbachelder/trap
After=network.target

[Service]
Type=simple
User=dan
WorkingDirectory=/home/dan/src/trap
EnvironmentFile=/home/dan/src/trap/.env.local
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/home/dan/.volta/tools/image/node/22.22.0/bin/node ./node_modules/next/dist/bin/next start -p 3002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trap-server

[Install]
WantedBy=default.target
