# systemd/trap-bridge.service
[Unit]
Description=Trap Chat Bridge
Documentation=https://github.com/dbachelder/trap
After=network.target trap-server.service
Wants=trap-server.service

[Service]
Type=simple
WorkingDirectory=/home/dan/src/trap
EnvironmentFile=/home/dan/src/trap/.env.local
ExecStart=/home/dan/.volta/tools/image/node/22.22.0/bin/node ./node_modules/tsx/dist/cli.cjs worker/chat-bridge.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trap-bridge

[Install]
WantedBy=default.target
